# scripts/services/startup.service
[Unit]
Description=UpdateGit

[Service]
ExecStart=~/gitrepos/scripts/startup.sh

[Install]
WantedBy=default.target 