# scripts/services/startup.service
[Unit]
Description=UpdateGit

[Service]
ExecStart=/usr/bin/startup.sh
Type=oneshot

[Install]
WantedBy=default.target 