# scripts/services/startup.service
[Unit]
Description=UpdateGit

[Service]
ExecStart=/usr/bin/startup.sh

[Install]
WantedBy=multi-user.target 